# supervisord.conf
[supervisord]
logfile = supervisord.log
logfile_maxbytes = 50MB
loglevel = info
nodaemon = true
identifier = supervisor
nocleanup = true
pidfile = supervisord.pid
directory= .

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = localhost:9001

[supervisorctl]
prompt = datadog

[program:agent]
command=python agent.py
stdout_logfile=collector.log
redirect_stderr=true
priority=999
startsecs=2

[program:dogstatsd]
command=python dogstatsd.py
stdout_logfile=dogstatsd.log
redirect_stderr=true
priority=998
startsecs=3

[group:datadog-agent]
programs=agent,dogstatsd
